# containers/tinyrss-cron/crontab
* * * * * root php /path/to/update/job/php
